# crontab
PATH=/usr/local/bin:$PATH

@daily root /app/manage.py update_exchange_rates >/proc/1/fd/1 2>/proc/1/fd/2
